# marketpulse-daemon.service
# MarketPulse Daemon Systemd Service
# /etc/systemd/system/marketpulse-daemon.service

[Unit]
Description=MarketPulse News Crawler Daemon
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=marketpulse
Group=marketpulse
WorkingDirectory=/opt/marketpulse
Environment="PATH=/opt/marketpulse/venv/bin"

# PostgreSQL 사용 예시
ExecStart=/opt/marketpulse/venv/bin/python daemon.py --db-url postgresql://marketpulse:password@localhost:5432/marketpulse

# SQLite 사용 예시 (주석 해제)
# ExecStart=/opt/marketpulse/venv/bin/python daemon.py

# 재시작 정책
Restart=always
RestartSec=10

# 로그
StandardOutput=journal
StandardError=journal
SyslogIdentifier=marketpulse-daemon

# 리소스 제한 (선택사항)
MemoryLimit=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
